# lib/systemd/system/bluetooth.service.d/10-experimental-mode.conf
[Service]
# Clear previous ExecStart
ExecStart=
ExecStart=/usr/libexec/bluetooth/bluetoothd --experimental
